# doc/qemu/rc.local
#!/bin/sh

mount -t 9p -o trans=virtio host0 /mnt
